# provisionning/roles/nodejs/templates/etc/profile.d/nodejs.sh.j2
export NODEJS_HOME={{ nodejs['home_dir'] }}
export PATH=$PATH:$NODEJS_HOME/bin