# provisionning/roles/zookeeper/templates/profile.d/zookeeper.sh.j2
export ZOOKEEPER_HOME={{ zookeeper['home_dir'] }}
export PATH=$PATH:$ZOOKEEPER_HOME/bin