# provisionning/roles/storm_common/templates/etc/profile.d/storm.sh.j2
export STORM_HOME={{ storm['home_dir'] }}
export PATH=$PATH:$STORM_HOME/bin
export STORM_ENV={{ storm['env'] }}