# provisionning/roles/redis/templates/etc/profile.d/redis.sh.j2
export REDIS_HOME={{ redis['home_dir'] }}
export PATH=$PATH:{{ redis['bin_dir'] }}
