# provisionning/roles/kafka/templates/etc/profile.d/kafka.sh.j2
export KAFKA_HOME={{ kafka['link'] }}
export PATH=$PATH:$KAFKA_HOME/bin