# provisionning/roles/java/templates/profile.d/java.sh.j2
export JAVA_HOME={{ java['home_dir'] }}
export PATH=$PATH:$JAVA_HOME/bin
